# crontab
0 5 */2 * * /a/do-backup >> /tmp/log 2>&1
