# crontab
0 5 * * * /a/do-backup >> /tmp/log 2>&1
