# crontab
0 5 * * * . /a/env && /a/do-backup >> /tmp/log 2>&1
