# crontab
0 6 * * * /a/do-backup >> /tmp/log 2>&1
